# asus-fancontrol.service
[Unit]
Description=ASUS UX32V fan control

[Service]
ExecStart=/usr/bin/asus-fancontrol
Restart=always

[Install]
WantedBy=graphical.target
