# asus-fancontrol.service
[Unit]
Description=ASUS UX31A fan control

[Service]
ExecStart=/usr/bin/asus-fancontrol
Restart=always

[Install]
WantedBy=graphical.target
